# rm-render@.service
[Unit]
Description=rm-render: watch rmfakecloud sync and render documents to PDF
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=%i
ExecStart=/bin/bash %h/src/rm-render/watch.sh %h/remotehomes/remarkable
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rm-render

[Install]
WantedBy=multi-user.target
