# rm-render@.service
[Unit]
Description=rm-render: watch rmfakecloud sync and render documents to PDF
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=%i
Environment="PATH=/home/%i/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/bin/bash /home/%i/src/rm-render/watch.sh /home/%i/remotehomes/remarkable
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rm-render

[Install]
WantedBy=multi-user.target
